# recipes-glowforge/glowforge-compat/glowforge-compat/glowforge-compat.list
/glowforge
/sv
/etc/glowforge.conf
/etc/init.d/factory-reset
/etc/init.d/head-update
/etc/init.d/runit
/usr/bin/blhost
/usr/bin/blhost-head-serial.sh
/usr/bin/bottle.py
/usr/bin/configure_wifi.sh
/usr/bin/create_virtual_if.sh
/usr/bin/doggo
/usr/bin/factory_reset.sh
/usr/bin/force_timeout
/usr/bin/fwup
/usr/bin/gen_ap_mac.py
/usr/bin/head-update.sh
/usr/bin/in_recovery_mode
/usr/bin/machine_hostname
/usr/bin/machine_name
/usr/bin/mmc_device
/usr/bin/mmc_partition
/usr/bin/process_crashdump.sh
/usr/bin/save_pstore.sh
/usr/bin/set_ap_name.sh
/usr/bin/watchdog.sh
/usr/bin/watchdog_keepalive_start.sh
/usr/bin/watchdog_keepalive_stop.sh
/usr/lib/libstdc\+\+.so.6
/usr/lib/libturbojpeg.so.0
/usr/lib/libuv.so.1
/usr/lib/libwebsockets.so.12
